# office365-collector.service
[Unit]
Description=Office365 Audit Log Collector
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/new-siem/office365-audit-log-collector
ExecStart=/home/ubuntu/new-siem/office365-audit-log-collector/target/release/office_audit_log_collector --config /home/ubuntu/new-siem/office365-audit-log-collector/config/config.yaml
Restart=always
RestartSec=10
LimitNOFILE=65535

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=office365-collector

[Install]
WantedBy=multi-user.target
